# rootfs-overlay/etc/systemd/system/openwhale.service
[Unit]
Description=OpenWhale AI Agent Platform
Documentation=https://github.com/viralcode/openwhale
After=network.target
Wants=network.target

[Service]
Type=simple
User=ainux
Group=ainux
WorkingDirectory=/opt/ainux/openwhale

ExecStart=/usr/bin/node openwhale.mjs

Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=120

# Environment
Environment=NODE_ENV=production
Environment=PORT=7777
Environment=AINUX_MODE=true
Environment=HOME=/home/ainux

# Allow OpenWhale to execute commands and access files
NoNewPrivileges=false
ProtectSystem=false
ReadWritePaths=/home/ainux /opt/ainux/openwhale /tmp /home/ainux/Works

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openwhale

[Install]
WantedBy=multi-user.target
